# update-blocklists.service
[Unit]
Description=Update gravity lists

[Service]
ExecStart=pihole -g
